# confs/supervisord/certs.conf
[program:certs]
command=/usr/bin/python /opt/certificate/certs.py
redirect_stderr=true
stdout_logfile=/var/log/certs.log
stdout_logfile_maxbytes=100MB
