# confs/supervisord/consumer.conf
[program:consumer]
command=/usr/bin/python /opt/certificate/consumer.py
redirect_stderr=true
stdout_logfile=/var/log/certs_consumer.log
stdout_logfile_maxbytes=100MB


